# files/etc/crontabs/root
# Crontab for OJS Schedule tasks
0 * * * *   ojs-run-scheduled
